# flights.cron
0 12 * * * /home/ubuntu/flights/main.sh
